# root/etc/init.d/generic-daemon
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="${SVCNAME} daemon"
pidfile="/run/${SVCNAME}.pid"
start_stop_daemon_args="--user ${user:-nobody}:${group:-nobody}"
